# config/systemd/homeserver-backup.service
[Unit]
Description=Immich backup
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/scripts/backup.sh
User=root
